# cex-scanner.service
[Unit]
Description=ArbitrageSmartAI High-Speed CEX Price Scanner Background Service
After=network.target

[Service]
Type=simple
User=HP
WorkingDirectory=/home/HP/arbitrage-smart-ai
ExecStart=/usr/bin/python3 scripts/arbitrage_live_engine.py

# Auto Recovery Settings
Restart=always
RestartSec=5

# Resource Quotas Hardening
CPUQuota=35%
MemoryMax=500M

# Security Sandboxing Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
