# dex-scanner.service
[Unit]
Description=ArbitrageSmartAI High-Speed DEX & Cross-Chain Scanner Service
After=network.target

[Service]
Type=simple
User=HP
WorkingDirectory=/home/HP/globalpulsenewsmedia
ExecStart=/usr/bin/python3 scripts/dex_live_scanner.py

# Auto Recovery Settings
Restart=always
RestartSec=5

# Resource Quotas Hardening
CPUQuota=35%
MemoryMax=500M

# Security Sandboxing Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
