# frontend-server.service
[Unit]
Description=ArbitrageSmartAI Static Frontend Server (Port 8080)
After=network.target

[Service]
Type=simple
User=HP
WorkingDirectory=/home/HP/globalpulsenewsmedia
ExecStart=/usr/bin/npx serve -l 8080 .
Restart=always
RestartSec=5

# Resource Quotas Hardening
CPUQuota=20%
MemoryMax=300M

# Security Sandboxing Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
