# news-scraper.service
[Unit]
Description=Continuous News Scraper Background Service
After=network.target

[Service]
Type=simple
User=HP
WorkingDirectory=/home/HP/globalpulsenewsmedia
ExecStart=/usr/bin/python3 run_scraper_loop.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
